# hooks/apply-wayland-flags.service
[Unit]
Description=Apply Wayland flags to Electron/Chromium desktop files

[Service]
Type=oneshot
ExecStart=@SCRIPT@ --full
