# hooks/apply-wayland-flags.path
[Unit]
Description=Watch for new or updated desktop files

[Path]
PathModified=/usr/share/applications/

[Install]
WantedBy=default.target
